# dotfiles/Ydotool/.config/autostart/ydotoold.desktop
[Desktop Entry]
Name=Ydotoold
Type=Application
Exec=systemctl --user start ydotoold
Terminal=false