# dotfiles/Kitty/.config/autostart/kitty-quick-terminal.desktop
[Desktop Entry]
Name=Kitty Quick Access Terminal
Type=Application
Exec=systemctl --user start kitty-quick-terminal
Terminal=false