# dotfiles/MyScripts/.config/autostart/user_dbus.desktop
[Desktop Entry]
Name=User DBus
Type=Application
Exec=systemctl --user start user_dbus.service
Terminal=false